# supervisord.conf
[supervisord]
nodaemon=true

#[program:fcgiwrap]
#command = /usr/sbin/fcgiwrap
#socket = unix:/var/run/fcgiwrap.socket


[fcgi-program:fcgiwrap]
command = /usr/sbin/fcgiwrap
user = www-data
socket = unix:///var/run/supervisor/%(program_name)s.sock
socket_owner = www-data:www-data
socket_mode = 0770
autorestart=true
autostart=true
startsecs=1
startretries=3
stopsignal=QUIT
stopwaitsecs=10
redirect_stderr=true
environment=PATH='/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin'

[program:nginx]
#command=/bin/bash -c "service fcgiwrap start && /usr/sbin/nginx -g 'daemon off;'"
command=/usr/sbin/nginx -g 'daemon off;'
#autorestart=true
